# .docker/php/secu-supervisor.ini
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/var/log/supervisor/supervisord.pid

[program:php-fpm]
command=php-fpm -F
process_name=%(program_name)s
autorestart=false
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:cron]
command=/usr/sbin/crond -f
process_name=%(program_name)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[group:secu]
programs=
    php-fpm,
    cron
autostart=true
